# systemd/.config/systemd/user/batteryNotif.service
[Unit]
Description=Low Battery Warning Notification

[Service]
ExecStart=%h/Linux-Information-Notification-Scripts/lowBatteryNotif.sh
